# packages/utah-project/utah_project-1.7.0-py3-none-any.whl/utah_project-1.7.0.data/data/server_config/centos8/etc/profile.d/utah_arches_profile.sh
export UTAH_APP_HOME=/opt/utah/arches
export UTAH_VARIABLE_DATA=/var/opt/utah/arches